# Services/logSSH.service
[Unit]
Description=SSH Monitor service

[Service]
Type=simple
ExecStart=/bin/bash /usr/bin/LOGssh.sh

[Install]
WantedBy=multi-user.target
